# dev/docker/etc/services.d/golang/run
#!/usr/bin/execlineb -P
with-contenv s6-setuidgid golang /app/gatekeeper --log-level=debug